# examples/turtlesim/s6-rc.d/40-turtlesim/run
#!/command/with-contenv bash

source /opt/ros/$(rosversion -d)/setup.bash
rosrun turtlesim turtlesim_node __ip:=$(hostname -i)