# examples/turtlesim/s6-rc.d/10-xvfb/run
#!/command/with-contenv bash

exec Xvfb -screen 0 500x500x24 > /dev/null 2>&1