# examples/turtlesim/s6-rc.d/50-notify/up
#!/usr/bin/with-contenv bash

bash -c "echo \"TurtleSim is Ready!
Access URL: http://127.0.0.1:8080/vnc.html
Password: $(cat /run/vncpass-raw)\" | boxes -a jc -d whirly -p h4t1"